# launchd/com.user.signal-scan.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN">
<plist version="1.0"><dict>
  <key>Label</key><string>com.user.signal-scan</string>
  <key>ProgramArguments</key><array>
    <string>/bin/bash</string><string>-lc</string><string>~/Automation/bin/notify_signals.sh</string>
  </array>
  <key>StartInterval</key><integer>3600</integer>
  <key>RunAtLoad</key><true/>
  <key>StandardOutPath</key><string>~/Automation/logs/signal_scan.out</string>
  <key>StandardErrorPath</key><string>~/Automation/logs/signal_scan.err</string>
</dict></plist>
